# r-cyber.service
[Unit]
Description=R-CYBER Patrimoine — Gunicorn WSGI server
After=network.target

[Service]
Type=notify

# ── Utilisateur / groupe ──────────────────────────────────────
User=rcyber
Group=www-data

# ── Répertoire de travail (racine du projet Django) ───────────
WorkingDirectory=/var/www/rcyber/app

# ── Variables d'environnement ─────────────────────────────────
EnvironmentFile=/var/www/rcyber/app/.env

# ── Commande de lancement ─────────────────────────────────────
# core.wsgi = app Django (adapte "core" si ton module WSGI est différent)
ExecStart=/var/www/rcyber/app/venv/bin/gunicorn \
    --workers 3 \
    --bind 127.0.0.1:8000 \
    --timeout 300 \
    --access-logfile /var/log/gunicorn/rcyber_access.log \
    --error-logfile  /var/log/gunicorn/rcyber_error.log \
    --log-level info \
    core.wsgi:application

# ── Redémarrage automatique en cas de crash ───────────────────
Restart=on-failure
RestartSec=5s

# ── Limites ───────────────────────────────────────────────────
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
